# www/crontab
* * * * * root /usr/local/bin/node /app/fetch_gifs.js >> /var/log/cron.log 2>&1